# scripts/server/systemd/spacers-dataset.service.example
[Unit]
Description=Spacers dataset runner
After=qwen-vllm.service
[Service]
Type=oneshot
EnvironmentFile=%h/.config/spacers-agent/server.env
WorkingDirectory=${PROJECT_ROOT}
ExecStart=${PROJECT_ROOT}/scripts/server/run_dataset.sh
